# ops/systemd/amnezia-nat-setup.service
[Unit]
Description=AmneziaWG NAT and route for VPN subnets (10.8.1.0/24, 10.66.66.0/24)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
# Give amnezia-awg a few seconds to bring up awg0 so route add works
ExecStartPre=/bin/sleep 5
ExecStart=/bin/bash /opt/vpn-suite/ops/amnezia-nat-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
