# infra/systemd/units/vpn-suite-backup-monitoring.timer
[Unit]
Description=VPN Suite - scheduled monitoring config backups (daily)

[Timer]
OnCalendar=*-*-* 04:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
